# supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/etc/init.d/nginx start
autostart=true

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
